# launchd/com.rajesh.pdfqdrant.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Search service (server.py). Starts at login, auto-restarts on crash,
  so http://localhost:8131/ is always available.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rajesh.pdfqdrant.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_ROOT__/venv/bin/python</string>
        <string>__PROJECT_ROOT__/server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_ROOT__</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>__PROJECT_ROOT__/data/server.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_ROOT__/data/server.log</string>
</dict>
</plist>
