# launchd/com.rajesh.pdfqdrant.indexer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily PDF indexer schedule. Runs indexer.py every day at 08:15
  (offset from the v1 app's 08:00 so both aren't hashing the same
  files simultaneously). Unchanged files are skipped cheaply.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rajesh.pdfqdrant.indexer</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_ROOT__/venv/bin/python</string>
        <string>__PROJECT_ROOT__/indexer.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_ROOT__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__PROJECT_ROOT__/data/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_ROOT__/data/launchd.log</string>
</dict>
</plist>
